# roles/keepalived_exporter/templates/prometheus-keepalived-exporter.service.j2
[Unit]
After=network-online.target

[Service]
User=root
ExecStart=/usr/local/bin/prometheus-keepalived-exporter

[Install]
WantedBy=multi-user.target
